# systemd/tomcat-start.service
[Unit]
Description = Tomcat start service
#After =
#Requires =

[Service]
Type = oneshot
ExecStart = /usr/bin/systemctl start tomcat

