# systemd/tomcat-stop.service
[Unit]
Description = Tomcat stop service
#After =
#Requires = 

[Service]
Type = oneshot
ExecStart = /usr/bin/systemctl stop tomcat

